# .config/systemd/user/sommelier-x@.service.d/override.conf
[Service]
Environment="SOMMELIER_FRAME_COLOR=#111213"
